# logging/elasticsearch/es-supervisord.conf
[program:elasticsearch]      ; 添加elasticsearch服务  
command = /opt/tomcat/elasticsearch/bin/elasticsearch ;启动程序路径可以带参数
autostart = true
startsecs = 5
autorestart = true
startretries = 3              ; 其他的就不写了，配置文件里也有英文解释
user = tomcat                 ; 以哪个用户启动
redirect_stderr = true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 20
stdout_logfile = /opt/tomcat/logs/elasticsearch_super.log
