# logging/elasticsearch/kibana-supervisord.conf
[program:kibana]
directory = /opt/tomcat/kibana/bin/
command = /opt/tomcat/kibana/bin/kibana
autostart = true
startsecs = 5
autorestart = true
startretries = 3
user = tomcat     
redirect_stderr = true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 20
stdout_logfile = /opt/tomcat/logs/kibana_super.log


[group:elk]                              ; 组管理非常方便，我们可以start|stop|status
programs=elasticsearch,logstash,kibana   ; elk:* 来启动或关闭所有此组的服务，非必须开启。
