# logging/fluent-bit/fluentb.conf
[program:fluentb]
directory=/work/opt/fluentb/build
command=/work/opt/fluentb/build/bin/fluent-bit -c /work/opt/fluentb/build/conf/fluent.conf
;process_name=%(process_num)02d
numprocs=1
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/data/logs/fluent/fluent.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
