# tomcat/supervisord.conf
[program:tomcat]
command=/opt/tomcat/bin/catalina.sh run
environment=JAVA_HOME="/usr/local/jdk/",JAVA_BIN="/usr/local/jdk/bin"
directory=/opt/tomcat
autostart = true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=10M
priority=32
numprocs=1
stopsignal=TERM
stopasgroup=true
startsecs=10
startretries=999
log_stdout=true
log_stderr=true
logfile=/var/log/tomcat/supervisord-tomcat.out
logfile_maxbytes=20MB
logfile_backups=10
