# docker-desktop/linuxserver/docker-baseimage-selkies-ubuntu2404/root/etc/s6-overlay/s6-rc.d/svc-dbus/run
#!/usr/bin/with-contenv bash

# Folder setup
mkdir -p /run/dbus
chown abc:abc /run/dbus
rm -f /run/dbus/pid

# Run dbus
exec s6-setuidgid abc \
  dbus-daemon \
    --system \
    --nofork \
    --nosyslog 
