# postgresql/supervisord.conf
[program:postgresql_server]
command=/bin/bash -c "exec /usr/local/postgresql_9.6.2/bin/pg_ctl start -l /usr/local/postgresql_9.6.2/data/log/pg_server.log
directory=/usr/local/postgresql_9.6.2/bin"
autostart=true
autorestart=true
user = postgres
