# docker-desktop/xrdp/docker-antix/rootfs/etc/sv/dbus/run
#!/bin/sh
exec 2>&1
exec dbus-daemon --system --nofork --nopidfile
