# tracing/pinpoint/pinpoint-agent-attach-example/nginx-php/ubuntu18.04-nginx-php7.2/conf/etc/supervisor.d/pinpoint-collector-agent.conf
[group:pinpoint]
programs=pinpoint-collector-agent
priority=20

[program:pinpoint-collector-agent]
directory = /opt/pinpoint-collector-agent/
command = python3 run.py
process_name=%(program_name)s
startsecs = 0
autostart = true
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
