# docker-desktop/xrdp/docker-ubuntu-lxde/xrdp/sv/xrdp-sesman/run
#!/bin/sh
exec 2>&1
rm -f /var/run/xrdp-sesman/sesman.sock
exec xrdp-sesman --nodaemon
