# nginx/nginx-php/php-fpm-supervisord.conf
[program:php-fpm] 
command=/usr/local/php/sbin/php-fpm --nodaemonize
autostart=true
autorestart=true
user=root 
numprocs=1 
stdout_logfile=/var/log/supervisor-php-fpm.log
