# nginx/nginx-php/php-app-supervisor.conf
[program:pmapp]
command=/usr/bin/php -f /var/www/MosquittoPhpApp/Mqtt_test.php
;process_name=%(process_num)02d 
;numprocs=5 #启动几个进程
autostart=true #随着supervisord的启动而启动
autorestart=true #自动启动
startsecs=1 #程序重启时候停留在runing状态的秒数
startretries=10 #启动失败时的最多重试次数
redirect_stderr=true #重定向stderr到stdout
stdout_logfile=/var/www/MosquittoPhpApp/pmapp.log #stdout文件
user=daemon
