# nginx/nginx-php/supervisord.conf
;关闭后台模式，在执行的命令后面添加-g "daemon off;"
[program:nginx]
command=/usr/local/tools/nginx/sbin/nginx -g 'daemon off;'
process_name=%(program_name)s
priority=36
autostart=true
startretries=3
autorestart=true
user=root
