# nginx/nginx-oss/nginx-loadbalancer/rootfs/etc/services.d/nginx-error-log/run
#!/usr/bin/execlineb -P
logutil-service -f /var/run/s6/nginx-error-log-fifo /var/log/nginx-error-logs
