# nginx/openresty/debian10-openresty-repack/after-install.sh
#!/bin/bash
ldconfig
#systemctl daemon-reload >/dev/null 2>&1 ||:
#systemctl enable nginx
#systemctl restart nginx
