# consul/supervisord.conf
[program:consul]
directory=/usr/local/bin/
command=/usr/local/bin/consul agent -config-dir=/etc/consul.d
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/consul/consul.log
