# docker-desktop/xrdp/docker-ubuntu-lxde/docker-antix/rootfs/etc/sv/rustdesk/run
#!/bin/sh
exec 2>&1

# 设置环境变量
export PULSE_LATENCY_MSEC=60
export PIPEWIRE_LATENCY=1024/48000

ulimit -n 100000
exec /usr/bin/rustdesk --service
