# monitoring/prometheus/node-exporter/node-exporter-supervisord.conf
[program:node_exporter]
command=/usr/local/bin/node_exporter
stdout_logfile=/tmp/prometheus/prometheus.log
autostart=true
autorestart=true
startsecs=5
priority=10
user=root
stopasgroup=true
killasgroup=true
