# bigdata/hbase/supervisord.conf
[program:regionserver]

command=bash /usr/hdp/current/hbase-regionserver/bin/hbase-daemon.sh --config /usr/hdp/current/hbase-regionserver/conf foreground_start regionserver

process_name=%(program_name)s

numprocs=1

;directory=/usr/hdp/current/hbase-regionserver/bin/

autostart=true

autorestart=true

startretries=3

exitcodes=0,2

stopsignal=INT

user=hbase

stdout_logfile=regionserver_out.log

stdout_logfile_maxbytes=1MB

stdout_logfile_backups=10

stdout_capture_maxbytes=1MB

stderr_logfile=test_err.log

stderr_logfile_maxbytes=1MB

stderr_logfile_backups=10

stderr_capture_maxbytes=1MB
