# monitoring/prometheus/nginx-exporter-supervisord.conf
[program:nginx_exporter]
command=/usr/local/prometheus/nginx-vts-exporter/bin/nginx-vts-exporter -nginx.scrape_uri=http://10.100.xx.xxx/status/format/json
stdout_logfile=/tmp/prometheus/nginx-vts-exporter.log
autostart=true
autorestart=true
startsecs=5
priority=1
user=root
stopasgroup=true
killasgroup=true

