# logging/logstash/logstash-supervisord.conf
[program:logstash]
directory = /opt/tomcat/logstash/bin/
command = /opt/tomcat/logstash/bin/logstash -f /opt/tomcat/testlog.conf
autostart = true
startsecs = 5
autorestart = true
startretries = 3
user = tomcat
redirect_stderr = true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 20
stdout_logfile = /opt/tomcat/logs/logstash_super.log
