# kafka/kafka/supervisord.conf
[program:kafka]
command=/opt/kafka/bin/kafka-server-start.sh /opt/kafka/config/server.properties
directory=/opt/kafka/
autostart=true
autorestart=true
priority=30
numprocs=1
stopsignal=TERM
stopasgroup=true
startsecs=10
startretries=999
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/supervisord-kafka.out
logfile_maxbytes=20MB
logfile_backups=10
