# baseimages/ubuntu/ubuntu-24.04-runit/sv/php-fpm/run
#!/bin/bash
exec /usr/sbin/php-fpm8.2 --fpm-config /etc/php/8.2/fpm/php-fpm.conf -R

