# docker-desktop/blackarch-novnc/supervisord.conf
[supervisord]
nodaemon=true
#loglevel=trace
logfile=/var/log/supervisor/supervisord.log

[inet_http_server]
port=127.0.0.1:9001

[program:sshkeys]
command=/usr/bin/ssh-keygen -A
autostart=true
autorestart=false
exitcodes=0
startsecs=0
priority=100

[program:sshdir]
command=mkdir -p /run/sshd
autostart=true
autorestart=false
exitcodes=0
startsecs=0
priority=1

[program:sshd]
command=/sbin/sshd -D
autostart=true
autorestart=true
restartpause=3

[program:X11]
command=/sbin/Xvfb :0 -screen 0 1280x800x24
autorestart=true
priority=300

[program:x11vnc]
command=/sbin/x11vnc -localhost -cursor arrow -nopw -display :0 -xkb -noxrecord -noxdamage -wait 5 -shared -forever
autorestart=true

[program:novnc]
command=/sbin/novnc --vnc localhost:5900 --listen 8080
priority=450
autorestart=true

[program:startxfce4]
priority=460
command=/sbin/dbus-launch startxfce4
autorestart=true
environment=DISPLAY=":0",HOME="/root",USER="root"
