# jenkins/kube-jenkins-dind-agent/jdk8/debian10-kube-jenkins-dind-agent-openjdk/rootfs/etc/services.d/dockerd/run
#!/bin/execlineb -P

sudo /usr/bin/dockerd
