# zookeeper/zookeeper/supervisord.conf
[program:zookeeper]
command=/opt/zookeeper/bin/zkServer.sh start-foreground
user=root
directory=/opt/zookeeper/
autostart=true
autorestart=true
priority=30
numprocs=1
stopsignal=TERM
stopasgroup=true
startsecs=10
startretries=31
log_stdout=true
log_stderr=true
logfile=/var/log/zookeeper/supervisord-zookeeper.out
logfile_maxbytes=20MB
logfile_backups=10
