# docker-desktop/xrdp/docker-ubuntu-lxde/xrdp/rootfs/etc/sv/xrdp-sesman/run
#!/bin/sh
exec 2>&1

rm -f /var/run/xrdp/xrdp-sesman.pid
exec xrdp-sesman --nodaemon
