# redis/supervisord.conf
;关闭redis的后台模式，也就是修改redis.conf中的daemon为no
[program:redisd]
command=/usr/local/redis/bin/redis-server
process_name=%(program_name)s
priority=35
autostart=true
startretries=3
autorestart=true
user=root
